# roles/cis/templates/inspec.timer
[Unit]
Description=Runs inspec and submits data to elasticsearch

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target
